# filebeat.conf
[supervisord]
nodaemon=true

[program:filebeat]
command = /bin/filebeat -c /tmp/filebeat.yml